# var/lib/dpkg/info/iputils-rdisc.list
/.
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/rdisc.service
/lib/systemd/system-preset
/lib/systemd/system-preset/98-iputils-rdisc.preset
/sbin
/sbin/rdisc
