# var/lib/dpkg/info/iputils-rarpd.list
/.
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/rarpd@.service
/sbin
/sbin/rarpd
