# var/lib/dpkg/info/m4projects-stm32mp1.prerm
#!/bin/sh
[ "$1" != "upgrade" ] || exit 0
set -e
if type systemctl >/dev/null 2>/dev/null; then
	if [ -z "$D" ]; then
		systemctl stop st-m4firmware-load.service

		systemctl disable st-m4firmware-load.service
	fi
fi
if type systemctl >/dev/null 2>/dev/null; then
	if [ -z "$D" ]; then
		systemctl stop st-m4firmware-load.service

		systemctl disable st-m4firmware-load.service
	fi
fi
