# var/lib/dpkg/info/mesa-megadriver.list
/.
/usr
/usr/share
/usr/share/drirc.d
/usr/share/drirc.d/00-mesa-defaults.conf
