# var/lib/dpkg/info/ppp.list
/.
/etc
/etc/chatscripts
/etc/chatscripts/pap
/etc/ppp
/etc/ppp/chap-secrets
/etc/ppp/ip-down
/etc/ppp/ip-down.d
/etc/ppp/ip-down.d/92removedns
/etc/ppp/ip-up
/etc/ppp/ip-up.d
/etc/ppp/ip-up.d/08setupdns
/etc/ppp/options
/etc/ppp/pap-secrets
/etc/ppp/peers
/etc/ppp/peers/provider
/etc/ppp/ppp_on_boot
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/ppp@.service
/usr
/usr/bin
/usr/bin/poff
/usr/bin/pon
/usr/sbin
/usr/sbin/chat
/usr/sbin/pppd
