# var/lib/dpkg/info/pulseaudio-server.list
/.
/etc
/etc/dbus-1
/etc/dbus-1/system.d
/etc/dbus-1/system.d/pulseaudio-system.conf
/etc/default
/etc/default/volatiles
/etc/default/volatiles/volatiles.04_pulse
/etc/profile.d
/etc/profile.d/pulse_profile.sh
/etc/pulse
/etc/pulse/daemon.conf
/etc/pulse/default.pa
/etc/pulse/system.pa
/etc/xdg
/etc/xdg/autostart
/etc/xdg/autostart/pulseaudio.desktop
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/pulseaudio.service
/lib/systemd/system/pulseaudio.socket
/lib/systemd/system-preset
/lib/systemd/system-preset/98-pulseaudio-server.preset
/lib/udev
/lib/udev/rules.d
/lib/udev/rules.d/90-pulseaudio.rules
/usr
/usr/bin
/usr/bin/pactl
/usr/bin/pulseaudio
/usr/bin/start-pulseaudio-x11
