# var/lib/dpkg/info/resize-helper.prerm
#!/bin/sh
[ "$1" != "upgrade" ] || exit 0
set -e
if type systemctl >/dev/null 2>/dev/null; then
	if [ -z "$D" ]; then
		systemctl stop resize-helper.service

		systemctl disable resize-helper.service
	fi
fi
if type systemctl >/dev/null 2>/dev/null; then
	if [ -z "$D" ]; then
		systemctl stop resize-helper.service

		systemctl disable resize-helper.service
	fi
fi
