# var/lib/dpkg/info/rpcbind.prerm
#!/bin/sh
[ "$1" != "upgrade" ] || exit 0
set -e
	update-alternatives --remove  rpcinfo /usr/bin/rpcinfo.rpcbind
if type systemctl >/dev/null 2>/dev/null; then
	if [ -z "$D" ]; then
		systemctl stop rpcbind.service rpcbind.socket

		systemctl disable rpcbind.service rpcbind.socket
	fi
fi
