# var/lib/dpkg/info/sysstat.list
/.
/etc
/etc/sysconfig
/etc/sysconfig/sysstat
/etc/sysconfig/sysstat.ioconf
/etc/tmpfiles.d
/etc/tmpfiles.d/sysstat.conf
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/sysstat.service
/usr
/usr/bin
/usr/bin/cifsiostat
/usr/bin/iostat
/usr/bin/mpstat
/usr/bin/pidstat
/usr/bin/sadf
/usr/bin/sar
/usr/bin/tapestat
/usr/libexec
/usr/libexec/sa
/usr/libexec/sa/sa1
/usr/libexec/sa/sa2
/usr/libexec/sa/sadc
