# var/lib/dpkg/info/usbip.list
/.
/etc
/etc/modprobe.d
/etc/modprobe.d/usbip.conf
/etc/udev
/etc/udev/rules.d
/etc/udev/rules.d/99-usb-usbip.rules
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/usbip.service
/lib/systemd/system-preset
/lib/systemd/system-preset/98-usbip.preset
/usr
/usr/lib
/usr/lib/libusbip.so.0.0.1
/usr/sbin
/usr/sbin/stm32mp-usbip-bind-unbind.sh
/usr/sbin/usbip
/usr/sbin/usbipd
/usr/share
/usr/share/hwdata
/usr/share/hwdata/usb.ids
/usr/lib/libusbip.so.0
