# var/lib/dpkg/info/weston-init.list
/.
/etc
/etc/default
/etc/default/weston
/etc/profile.d
/etc/profile.d/weston_profile.sh
/etc/udev
/etc/udev/rules.d
/etc/udev/rules.d/71-weston-drm.rules
/etc/xdg
/etc/xdg/weston
/etc/xdg/weston/weston.ini
/home
/home/root
/home/root/README-CHECK-GPU
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/weston@.service
/lib/systemd/system-preset
/lib/systemd/system-preset/98-weston-init.preset
/sbin
/sbin/check-gpu
/usr
/usr/bin
/usr/bin/weston-start
/usr/sbin
/usr/sbin/weston.sh
/usr/share
/usr/share/weston
/usr/share/weston/backgrounds
/usr/share/weston/backgrounds/ST13028_Linux_picto_11_1366x768.png
/usr/share/weston/backgrounds/ST_1366x768.png
/usr/share/weston/icon
/usr/share/weston/icon/ST13345_Products_light_blue_24x24.png
/usr/share/weston/icon/space.png
/usr/share/weston/icon/utilities-terminal.png
