# etc/udev/rules.d/85-m4-dump.rules
SUBSYSTEM=="devcoredump", ACTION=="add", RUN+="/usr/sbin/stm32mp-m4fwdump.sh"
